# security/pfSense-pkg-softether5/pkg-plist
etc/rc.d/vpnserver.sh
pkg/softether5.xml
www/vpn_softether5_vpnserver.php
%%DATADIR%%/info.xml
